# keepalive.service
[Service]
ExecStart=/home/ssm-user/imperiumhashbot/keepalive.sh